# rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bashio
# Run the PHP-FPM daemon
bashio::log.info "Starting PHP-FPM..."

exec php-fpm7 --nodaemonize
